# com.harsh.claude-slack-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.harsh.claude-slack-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/harshkhandelwal/Desktop/MISC/slack-claude-bot/.venv/bin/python</string>
        <string>/Users/harshkhandelwal/Desktop/MISC/slack-claude-bot/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/harshkhandelwal/Desktop/MISC/slack-claude-bot</string>

    <!-- ensure `claude` CLI on PATH for the Agent SDK -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/harshkhandelwal/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <!-- no time ceiling on background tasks: long research runs to completion -->
        <key>CLAUDE_CODE_PRINT_BG_WAIT_CEILING_MS</key>
        <string>0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- bot.log is owned by the app's rotating logger; launchd captures only
         crash/startup output here so the two never fight over one file -->
    <key>StandardOutPath</key>
    <string>/Users/harshkhandelwal/Desktop/MISC/slack-claude-bot/bot.startup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/harshkhandelwal/Desktop/MISC/slack-claude-bot/bot.startup.log</string>
</dict>
</plist>
